# debian/sidedoor.upstart
description "SSH connection daemon"

start on (local-filesystems and net-device-up IFACE!=lo)

respawn

setuid sidedoor

script
    [ -f /etc/default/sidedoor ] && . /etc/default/sidedoor
    exec /usr/bin/sidedoor $OPTIONS $REMOTE_SERVER
end script
